# autoscan.online/devops/supervisor/conf.d/auto24-get-new.conf
[program:auto24-get-new]
user = www-data
command = /var/www/lk.autoscan.online/yii auto24/get-new
directory = /var/www/lk.autoscan.online
autostart=true
autorestart=true
numprocs=1
redirect_stderr=true
stdout_logfile = /var/www/lk.autoscan.online/log/auto24-get-new.log