# autoscan.online/devops/supervisor/conf.d/ss-get-new.conf
[program:ss-get-new]
user = www-data
command = /var/www/lk.autoscan.online/yii ss/get-new
directory = /var/www/lk.autoscan.online
autostart=true
autorestart=true
numprocs=1
redirect_stderr=true
stdout_logfile = /var/www/lk.autoscan.online/log/ss-get-new.log